# ops/launchd/com.andes.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.andes.api</string>
  <key>WorkingDirectory</key>
  <string>/Users/charlie/Projects/forks/andes-app/v2/backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>uv</string>
    <string>run</string>
    <string>andes-api</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/andes-api.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/andes-api.err</string>
</dict>
</plist>
